# OS/airootfs/usr/lib/systemd/system/defaultwall.service
[Unit]
Description=Change default wallpaper for cscOS

[Service]
ExecStart=/root/.wallpaper/defaultwall.py

[Install]
WantedBy=multi-user.target
